# pi-controller/empire-pi.service
[Unit]
Description=NoxuOS AI Empire Pi Controller
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/NoxuOS/pi-controller
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/pi/NoxuOS/pi-controller/pi-controller.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
